# wedding-website.service
[Unit]
Description=Wedding Website Nginx Server
After=network.target nginx.service
Wants=network.target
Requires=nginx.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecReload=/usr/sbin/nginx -s reload

# This service ensures the wedding website is properly configured with nginx
# The actual serving is handled by nginx itself

[Install]
WantedBy=multi-user.target
